# mysql/mysql
#!/lib/init/init-d-script
### BEGIN INIT INFO
# Provides:          mysqld
# Required-Start:    $syslog $time $remote_fs
# Required-Stop:     $syslog $time $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      
# Short-Description: mysql server
# Description:       Debian init script to start the daemon
#                    mysqld.
### END INIT INFO
DAEMON=/usr/sbin/mysqld
case "$1" in
  start)
    mkdir -p /var/run/mysqld
	chmod 777 -R /var/run/mysqld
	chown mysql:mysql -R /var/run/mysqld
	start-stop-daemon --start --quiet --exec /usr/sbin/mysqld -- --user=root
	;;
  *)
	exit 1
esac
